# data/sources/meta-densitron/recipes-support/hwclock-sync/files/hwclock-sync.service
[Unit]
Description=Sync RTC from system time after NTP

[Service]
Type=oneshot
ExecStart=/sbin/hwclock -f /dev/rtc1 --systohc --update

